# script[Unit] Description=chathistoryadminwhite Documentation=http:/www.chathistoryadminwhite.io/After=network.target  [Service] Type=simple ExecStart=/bin/python3.6/opt/chat_history/chathistoryadminwhite.service
[Unit]
Description=chathistoryadminwhite
Documentation=http://www.chathistoryadminwhite.io/
After=network.target

[Service]
Type=simple
ExecStart=/bin/python3.6 /opt/chat_history/add_chat_history_admin_ip.py
Restart=always

[Install]
WantedBy=multi-user.target
